# units/vault.target
[Unit]
Description=Vault is online and unsealed

[Install]
WantedBy=multi-user.target
